# .config/systemd/user/urxvtd.service
[Unit]
Description = Urxvt Daemon

[Service]
ExecStart = /usr/bin/urxvtd
Restart = always

[Install]
WantedBy = default.target
